# crontab
# Run scraper daily at 3:00 AM Beijing time
0 3 * * * cd /app && python scheduler.py >> /var/log/cron.log 2>&1

# Empty line required at end of cron file